# init
#!/sbin/openrc-run



#command=/bin/gikeud
command=/data/01/src/gikeud/main.py
command_args="-d -p 1d"
pidfile="/var/run/${RC_SVCNAME}.pid"
name="gikeud daemon"
command_background=true
 
description="Gikeud is a daemon that auto git pull and guess kernel branch available"

depend() {
	# Need net to sync
	need net
	after logger
}
